# deploy/rockbaybooty.service
[Unit]
Description=RockBayBooty Next.js app
Documentation=https://github.com/deenyx/rockbaybooty
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/rockbaybooty

# Load production env file. Create this file (mode 640) by hand on the VPS;
# it must not be committed to the repository.
EnvironmentFile=/var/www/rockbaybooty/.env.production
Environment=NODE_ENV=production
Environment=PORT=3000

ExecStart=/usr/bin/npx next start -p 3000

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rockbaybooty

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
